# activiti/files/supervisor.conf
[supervisord]
nodaemon=true

[program:activiti]
